# cronjob-logrotate.service
[Unit]
Description=logrotate cron job

[Service]
ExecStart=/usr/lib/logrotate

# vi: ft=systemd
