# cronjob-logrotate.timer
[Unit]
Description=logrotate cron job timer

[Timer]
OnCalendar=*-*-* 05:02:00
Persistent=true

[Install]
WantedBy=cronjobs.target

# vi: ft=systemd
